# crontab
#
# $Id: crontab,v 1.1 2000/08/19 00:59:17 statproc Exp $
#
# Obviously, this file has no direct bearing on the cronned tasks.
# It just serves as a reference to the live crontab for the statproc
# user.  It's the responsibility of people editing statproc's
# crontab to keep this file up-to-date
#

5       *       *       *       *       cd /usr/home/statproc/dfmon; ./dfmon dfmon.conf
#22     *       *       *       *       /usr/home/statproc/stats-rc5/hourly.pl 2>&1 > /usr/home/statproc/log/hourly.log
25      0       *       *       *       /usr/home/statproc/statsrun/daily.sh > /usr/home/statproc/log/daily.log
4       *       *       *       *       cd /usr/home/statproc/stats-proc/hourly && ./hourly.pl > /usr/home/statproc/log/lastrun.hourly 2> /usr/home/statproc/log/lastrun.hourly.err

